# crontab
0 * * * * ash /home/root/update-sim.sh /home/root
